# systemd/lescraper-morning.timer
[Unit]
Description=LeScraper Morning Run Timer (7am-12pm randomized)
Documentation=https://github.com/sam-packer/LeScraper

[Timer]
Unit=lescraper.service

# Run every day at 7:00 AM
OnCalendar=*-*-* 07:00:00

# Add 0-5 hours of randomization (7am-12pm window)
RandomizedDelaySec=18000

# If the system was off when the timer should have triggered, run it on next boot
Persistent=true

[Install]
WantedBy=timers.target
